# hassio-mcz-maestro-dev/rootfs/etc/services.d/mcz-maestro/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start service mcz maestro
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "Start project mcz-maestro"

cp /data/* /share

## copy config file to app with yaml extension
cat /data/options.json | yq -P > /app/dist/config.yaml

## Run your program
exec cd /app && node dist/main.js
